# github-action-executor.service
[Unit]
Description=GitHub Action Executor
After=network.target

[Service]
Type=simple
User=YOUR_USERNAME
WorkingDirectory=/path/to/project
Environment="PATH=/path/to/project/venv/bin:/usr/local/bin:/usr/bin:/bin"
EnvironmentFile=/path/to/project/.env
ExecStart=/path/to/project/venv/bin/uvicorn app:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
